# files/etc/hotplug.d/ieee80211/11_fix_wifi_mac_n30pro
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}
[ -n "$PHYNBR" ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

[ "$(board_name)" = "netcore,n30-pro" ] || exit 0

# The factory base address is transformed in the same way as the known-good
# ImmortalWrt N30 Pro adaptation, yielding stable and distinct radio MACs.
addr=$(mtd_get_mac_binary "Factory" 0x4)
[ -n "$addr" ] || exit 0

addr=$(macaddr_unsetbit "$(macaddr_setbit_la "$addr")" 25)
addr_2g=$(macaddr_unsetbit "$(macaddr_unsetbit "$addr" 26)" 27)
addr_5g=$(macaddr_setbit "$(macaddr_setbit "$addr" 26)" 27)

[ "$PHYNBR" = "0" ] && printf '%s\n' "$addr_2g" > "/sys${DEVPATH}/macaddress"
[ "$PHYNBR" = "1" ] && printf '%s\n' "$addr_5g" > "/sys${DEVPATH}/macaddress"

exit 0
